# ops/launchd/com.forgecapital.weekly-digest.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.forgecapital.weekly-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/tristanfischer/Developer/forge-capital-app/scripts/weekly-digest-cron.mjs</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tristanfischer/Developer/forge-capital-app</string>

    <!-- Monday 07:00 BST. launchd reads StartCalendarInterval in the
         user's local time zone, so on a machine set to Europe/London
         this fires at 07:00 BST (summer) / 07:00 GMT (winter) without
         any extra tz gymnastics. Weekday=1 is Monday per launchd's
         Sun=0 / Mon=1 / ... / Sat=6 convention. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/private/tmp/fc-weekly-digest.log</string>

    <key>StandardErrorPath</key>
    <string>/private/tmp/fc-weekly-digest.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NEXT_PUBLIC_SUPABASE_URL</key>
        <string><!-- PASTE YOUR NEXT_PUBLIC_SUPABASE_URL HERE --></string>

        <key>SUPABASE_SERVICE_ROLE_KEY</key>
        <string><!-- PASTE YOUR SUPABASE_SERVICE_ROLE_KEY HERE --></string>

        <key>GMAIL_CLIENT_ID</key>
        <string><!-- PASTE YOUR GMAIL_CLIENT_ID HERE --></string>

        <key>GMAIL_CLIENT_SECRET</key>
        <string><!-- PASTE YOUR GMAIL_CLIENT_SECRET HERE --></string>

        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
